# examples/thermal_camera/sdcard/etc/rc.local
#!/bin/sh

cd root/dist/thermal_camera_release && ./thermal_camera

